# docker/systemd/chatbotchefeu.timer
[Unit]
Description=ChatBotChef EU daily reminders

[Timer]
OnCalendar=daily
Persistent=true
Unit=chatbotchefeu-housekeeping.service

[Install]
WantedBy=timers.target
